# config_files/ishiki-bootstrap.service
# /etc/systemd/system/ishiki-bootstrap.service

[Unit]
Description=Ishiki Bootstrap Service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/ishiki/bootstrap/start.sh
RemainAfterExit=true
ExecStop=/opt/ishiki/bootstrap/stop.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target
